# docker_config/mysql/supervisord.conf
[supervisord]
nodaemon=false

[program:mysql]
command=/run.sh

[program:inspeqtor]
command=/usr/bin/inspeqtor -c /etc/inspeqtor -s /var/run/inspeqtor.sock
